# deploy/systemd/ai-host-observability.service
[Unit]
Description=AI Host Observability textfile collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=OUT_DIR=/var/lib/node_exporter/textfile_collector
ExecStart=/opt/ai-host-observability/collect-all.sh

[Install]
WantedBy=multi-user.target
